# .docker/php/crontab
* * * * * "/home/wwwroot/sf4/bin/console app:process-transactions"
# An empty line is required at the end of this file for a valid cron file.